# coolify-manager.service
[Unit]
Description=Coolify Multi-Agent Management System
After=network.target

[Service]
Type=simple
User=coolify
WorkingDirectory=/opt/coolify-manager
ExecStart=/usr/bin/python3 /opt/coolify-manager/run.py
Restart=always
RestartSec=10

# Environment variables
EnvironmentFile=/opt/coolify-manager/.env

# Logging
StandardOutput=append:/opt/coolify-manager/logs/stdout.log
StandardError=append:/opt/coolify-manager/logs/stderr.log

[Install]
WantedBy=multi-user.target
